# generator_control.service
[Unit]
Description=Powermate PM9400E Generator Control (Flask + GPIOZero)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/generatorpi
ExecStart=/usr/bin/python3 /home/pi/generatorpi/generator_control.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
